# ci/files/etc/cron.d/fix-cpu-scaling
@reboot root /usr/local/bin/fix-cpu-scaling.sh