# ci/files/etc/cron.d/take-snapshot
@daily root /usr/local/bin/snapshot.sh >/dev/null
