# configs/files/etc/cron.d/take-snapshot
0 7 * * * root /usr/local/bin/snapshot.sh >/dev/null
